# cityofhope-20131204-2000-public/etc/cron.d/tinymux
0 */4 * * *     tinymux     /opt/tinymux/BackupAll
